# crontab
* * * * * sh /home/pi/camera/capture.sh
@reboot cd /home/pi/camera/ && node index.js

